# edge/systemd/csi-aggregator.service
[Unit]
Description=CSI Sensing Edge Pipeline — UDP aggregator + DSP + GPU forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%DEPLOY_USER%
Group=%DEPLOY_USER%
WorkingDirectory=%DEPLOY_REPO%
Environment=PYTHONPATH=%DEPLOY_REPO%
ExecStart=%DEPLOY_REPO%/.venv/bin/python -m edge.aggregator \
    --config %DEPLOY_REPO%/edge/config/pipeline.yaml \
    --log-level INFO

# Performance tuning
CPUAffinity=2 3
Nice=-5

Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csi-aggregator

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%DEPLOY_REPO%

[Install]
WantedBy=multi-user.target
